# content/etc/s6-overlay/s6-rc.d/healthcheck/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

python3 /usr/local/bin/wg_healthcheck.py
